# deploy/ladder-proxy.service.d/cors.conf
# Drop-in for ladder-proxy.service. Install at:
#   /etc/systemd/system/ladder-proxy.service.d/cors.conf
# then: sudo systemctl daemon-reload && sudo systemctl restart ladder-proxy
#
# Set ALLOWED_ORIGINS to every origin that loads the engine, playground,
# explorer, or any other page that calls /api/ladder/*. Comma-separated,
# no spaces around commas. The ladder_proxy.py default is bitcoinghost.org
# only and will reject browser requests from ladder-script.org.

[Service]
Environment=ALLOWED_ORIGINS=https://ladder-script.org,https://www.ladder-script.org,http://localhost:8080,http://127.0.0.1:8080
